# playwright-duunitori.service
[Unit]
Description=Playwright Crawler - Duunitori.fi job listings

[Service]
Type=oneshot
ExecStart=/usr/bin/npx ts-node src/duunitori_crawler.ts
WorkingDirectory=/home/mikko/github/playwright-crawlers
EnvironmentFile=/opt/playwright-crawlers/.env

[Install]
WantedBy=default.target
