# systemd/services/playwright-fillaritori.service
[Unit]
Description=Playwright Crawler - Fillaritori.fi bicycle marketplace

[Service]
Type=oneshot
ExecStart=/usr/bin/npx ts-node src/fillaritori_crawler.ts
WorkingDirectory=/home/mikko/github/playwright-crawlers
EnvironmentFile=/opt/playwright-crawlers/.env

[Install]
WantedBy=default.target
