# systemd/services/playwright-mol.service
[Unit]
Description=Playwright Crawler - Tyomarkkinatori.fi (MOL) job listings

[Service]
Type=oneshot
ExecStart=/usr/bin/npx ts-node src/mol_crawler.ts
WorkingDirectory=/home/mikko/github/playwright-crawlers
EnvironmentFile=/opt/playwright-crawlers/.env

[Install]
WantedBy=default.target
