# systemd/services/playwright-krapinpaja.service
[Unit]
Description=Playwright Crawler - Krapinpaja.fi

[Service]
Type=oneshot
ExecStart=/usr/bin/npx ts-node src/krapinpaja_crawler.ts
WorkingDirectory=/home/mikko/github/playwright-crawlers
EnvironmentFile=/opt/playwright-crawlers/.env

[Install]
WantedBy=default.target
