# systemd/timers/playwright-krapinpaja-crawler.timer
[Unit]
Description=Run Krapinpaja crawler hourly
Documentation=https://github.com/frogshead/playwright-crawlers
Requires=playwright-crawler@krapinpaja.service

[Timer]
# Run every hour at minute 40
OnCalendar=*:40:00
Persistent=true
RandomizedDelaySec=60
# Explicitly specify which service to trigger
Unit=playwright-crawler@krapinpaja.service

[Install]
WantedBy=timers.target