# systemd/timers/playwright-fillaritori-crawler.timer
[Unit]
Description=Run Fillaritori crawler hourly
Documentation=https://github.com/frogshead/playwright-crawlers
Requires=playwright-crawler@fillaritori.service

[Timer]
# Run every hour at minute 20
OnCalendar=*:20:00
Persistent=true
RandomizedDelaySec=60
# Explicitly specify which service to trigger
Unit=playwright-crawler@fillaritori.service

[Install]
WantedBy=timers.target