# systemd/timers/playwright-tori-crawler.timer
[Unit]
Description=Run Tori crawler hourly
Documentation=https://github.com/frogshead/playwright-crawlers
Requires=playwright-crawler@tori.service

[Timer]
# Run every hour at minute 0
OnCalendar=*:00:00
# Enable persistent timers (catch up missed runs)
Persistent=true
# Add some randomization to prevent all crawlers from starting simultaneously
RandomizedDelaySec=60

[Install]
WantedBy=timers.target