# systemd/timers/playwright-mol-crawler.timer
[Unit]
Description=Run MOL crawler hourly
Documentation=https://github.com/frogshead/playwright-crawlers
Requires=playwright-crawler@mol.service

[Timer]
# Run every hour at minute 10
OnCalendar=*:10:00
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target