# systemd/deprecated/playwright-theseus-crawler.timer
[Unit]
Description=Run Theseus crawler hourly
Documentation=https://github.com/frogshead/playwright-crawlers
Requires=playwright-crawler@theseus.service

[Timer]
# Run every hour at minute 50
OnCalendar=*:50:00
Persistent=true
RandomizedDelaySec=60
# Explicitly specify which service to trigger
Unit=playwright-crawler@theseus.service

[Install]
WantedBy=timers.target