# systemd/deprecated/playwright-tavastia-crawler.timer
[Unit]
Description=Run Tavastia crawler hourly
Documentation=https://github.com/frogshead/playwright-crawlers
Requires=playwright-crawler@tavastia.service

[Timer]
# Run every hour at minute 30
OnCalendar=*:30:00
Persistent=true
RandomizedDelaySec=60
# Explicitly specify which service to trigger
Unit=playwright-crawler@tavastia.service

[Install]
WantedBy=timers.target